# files/system/usr/lib/systemd/user/topgrade.timer
[Unit]
Description=Delay Topgrade update by 5mn after login

[Timer]
OnBootSec=5m
Persistent=true

[Install]
WantedBy=timers.target
